# ssh_config/crontab
# Persist the github.com SSH connection, YOLO.
*/4 * * * * ssh -T git@github.com
